# services/xenbridge.service
[Unit]
Description=XenBridge SaaS API
After=network.target

[Service]
User=root
WorkingDirectory=/root/xenbridge/xenbridge
Environment="OPENAI_API_KEY=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"
ExecStart=/root/xenbridge/xenbridge/venv/bin/gunicorn -w 4 -b 0.0.0.0:8000 app:app
Restart=always

[Install]
WantedBy=multi-user.target